# files/etc/init.d/sd-wrt-policy
#!/bin/sh /etc/rc.common
USE_PROCD=1
START=96
STOP=4
start_service() { :; }
