# files/etc/init.d/sd-wrt-probe
#!/bin/sh /etc/rc.common
USE_PROCD=1
START=95
STOP=5
start_service() { :; }
