# files/etc/hotplug.d/iface/25-sd-wrt
#!/bin/sh
# Triggered by netifd on interface events.
# $INTERFACE = netifd logical name, $IFNAME = kernel device name.

[ "$ACTION" = ifup ] || [ "$ACTION" = ifdown ] || exit 0

. /lib/functions.sh

RUNDIR=/var/run/sd-wrt
mkdir -p "$RUNDIR"

config_load sd-wrt

# Find which sd-wrt section owns this kernel interface
matched_section=""
matched_role=""
_find() {
    local s="$1" ifname role
    config_get ifname "$s" ifname
    config_get role   "$s" role
    [ "$ifname" = "$IFNAME" ] || return
    matched_section="$s"
    matched_role="$role"
}
config_foreach _find link

[ -n "$matched_section" ] || exit 0

if [ "$ACTION" = ifup ]; then
    # Capture gateway before removing the default from main — policy's
    # setup_routing_tables won't find it after we delete the route.
    gw=$(ip route show dev "$IFNAME" | awk '/^default/{print $3; exit}')
    [ -n "$gw" ] && echo "$gw" > "$RUNDIR/$matched_section.gw"

    # Remove default route from main table for non-default-role interfaces.
    # sd-wrt-policy manages dedicated tables (101/102); without this every
    # interactive link that comes up hijacks all traffic in the main table.
    [ "$matched_role" != "default" ] && ip route del default dev "$IFNAME" 2>/dev/null

    # Reset state so probe re-evaluates on next cycle
    printf 'state=UNKNOWN\nlatency=0\nloss=0\nchanged=%d\n' "$(date +%s)" \
        > "$RUNDIR/$matched_section.state.tmp" \
        && mv "$RUNDIR/$matched_section.state.tmp" "$RUNDIR/$matched_section.state"
    : > "$RUNDIR/$matched_section.bad_count"
    : > "$RUNDIR/$matched_section.good_since"
else
    # ifdown: mark link DOWN immediately rather than waiting for probe timeout
    printf 'state=DOWN\nlatency=9999\nloss=100\nchanged=%d\n' "$(date +%s)" \
        > "$RUNDIR/$matched_section.state.tmp" \
        && mv "$RUNDIR/$matched_section.state.tmp" "$RUNDIR/$matched_section.state"
fi

# Signal probe for immediate re-evaluation
probe_pid=$(pgrep -f sd-wrt-probe | head -1)
[ -n "$probe_pid" ] && kill -USR1 "$probe_pid" 2>/dev/null

# Signal policy to rebuild routing tables and fwmark rules
policy_pid=$(cat /var/run/sd-wrt-policy.pid 2>/dev/null)
[ -n "$policy_pid" ] && kill -HUP "$policy_pid" 2>/dev/null
